# roles/mongodb/templates/thp.service.j2
[Unit]
Description={{ description }}
DefaultDependencies=no
After=sysinit.target local-fs.target
Before=mongod.service

[Service]
Type=oneshot
ExecStart={{ command }}

[Install]
WantedBy=basic.target