# roles/platform/templates/automation-platform.service.j2
[Unit]
Description=Itential Automation Platform Service
After=network.target
Requires=network.target

[Service]
Type=simple
User={{ platform_user }}
WorkingDirectory={{ platform_install_dir }}/current
ExecStart=/usr/bin/node --max-old-space-size=8192 server.js
Restart=on-failure
ExecReload=/bin/kill -USR2 $MAINPID
KillMode=mixed

[Install]
WantedBy=multi-user.target
